# etc/etc-nfsserver-rc/etc-nfsserver-rc.list
etc/rc.d/nfslocking
etc/rc.d/nfsd
etc/rc.d/mountd
