# xetc/etc-x11-rc/etc-x11-rc.list
etc/rc.d/xfs
etc/rc.d/xfs
etc/rc.d/xdm
etc/rc.d/xdm
etc/rc.d/fccache
etc/rc.d/fccache
