# etc/etc-nis-rc/etc-nis-rc.list
etc/rc.d/ypserv
etc/rc.d/yppasswdd
etc/rc.d/ypbind
