# etc/etc-sys-rc/tmp.list
etc/rc.subr
etc/rc.shutdown
etc/rc.local
etc/rc.d/wsmoused
etc/rc.d/wscons
etc/rc.d/virecover
etc/rc.d/veriexec
etc/rc.d/ttys
etc/rc.d/tpctl
etc/rc.d/syslogd
etc/rc.d/sysdb
etc/rc.d/sysctl
etc/rc.d/swap2
etc/rc.d/swap1
etc/rc.d/staticroute
etc/rc.d/securelevel
etc/rc.d/screenblank
etc/rc.d/savecore
etc/rc.d/rtclocaltime
etc/rc.d/root
etc/rc.d/rndctl
etc/rc.d/resize_root
etc/rc.d/random_seed
etc/rc.d/raidframeparity
etc/rc.d/raidframe
etc/rc.d/quota
etc/rc.d/pwcheck
etc/rc.d/perusertmp
etc/rc.d/newsyslog
etc/rc.d/mountcritremote
etc/rc.d/mountcritlocal
etc/rc.d/mountall
etc/rc.d/motd
etc/rc.d/modules
etc/rc.d/lvm
etc/rc.d/local
etc/rc.d/ldconfig
etc/rc.d/irdaattach
etc/rc.d/gpio
etc/rc.d/fsck_root
etc/rc.d/fsck
etc/rc.d/envsys
etc/rc.d/downinterfaces
etc/rc.d/dmesg
etc/rc.d/devpubd
etc/rc.d/cleartmp
etc/rc.d/cgd
etc/rc.d/ccd
etc/rc.d/bootconf.sh
etc/rc.d/bluetooth
etc/rc.d/blacklistd
etc/rc.d/apmd
etc/rc.d/accounting
etc/rc.d/SERVERS
etc/rc.d/NETWORKING
etc/rc.d/LOGIN
etc/rc.d/DISKS
etc/rc.d/DAEMON
etc/rc.conf
etc/rc
