# etc/etc-router-rc/tmp.list
etc/rc.d/routed
etc/rc.d/route6d
etc/rc.d/mrouted
etc/rc.d/ldpd
etc/rc.d/altqd
