# etc/etc-sysutil-rc/etc-sysutil-rc.list
etc/rc.d/wdogctl
etc/rc.d/powerd
etc/rc.d/moused
