# etc/etc-obsolete/etc-obsolete.list
etc/spamd.conf
etc/release-info
etc/rc.wscons
etc/rc.lkm
etc/rc.d/ypset
etc/rc.d/xntpd
etc/rc.d/systemfs
etc/rc.d/sdpd
etc/rc.d/poffd
etc/rc.d/lkm3
etc/rc.d/lkm2
etc/rc.d/lkm1
etc/rc.d/kerberos
etc/rc.d/gated
etc/rc.d/fsck.sh
etc/rc.d/fixsb
etc/rc.d/btuartd
etc/rc.d/bthcid
etc/rc.d/btdevctl
etc/rc.d/btcontrol
etc/rc.d/btconfig
etc/rc.d/btattach
etc/rc.d/NETWORK
etc/primes
etc/postinstall
etc/postfix/postfix-script
etc/postfix/postfix-files
etc/postfix/post-install
etc/obsolete/text
etc/obsolete/misc
etc/obsolete/man
etc/obsolete/games
etc/obsolete/etc
etc/obsolete/comp
etc/obsolete/base
etc/lkm.conf
etc/fstab.wd
etc/fstab.tmp
etc/fstab.sd7
etc/fstab.sd6
etc/fstab.sd5
etc/fstab.sd4
etc/fstab.sd3
etc/fstab.sd2
etc/fstab.sd1
etc/fstab.sd0
etc/fstab.sd
etc/fstab.nfs
etc/dhclient-script
etc/default/weekly.conf
etc/default/security.conf
etc/default/rc.conf
etc/default/monthly.conf
etc/default/daily.conf
etc/default
etc/bluetooth/btuartd.conf
etc/bluetooth/bluetooth.conf
etc/audit-packages.conf
dev/MAKEDEV.subr
